# profile/airootfs/etc/systemd/system/lga-learning@.service
[Unit]
Description=LGA constrained learning job %I
Documentation=file:///usr/share/doc/workspace/OPERATIONS.md
ConditionPathExists=/var/lib/lga/jobs/%I/job.json
After=local-fs.target lga-egressd.service

[Service]
Type=oneshot
User=lga-runner
Group=lga-ipc
ExecStart=/usr/bin/python /usr/local/lib/workspace/job_runner.py --policy /etc/lga/policy.toml --job-id %I
WorkingDirectory=/var/lib/lga/workspaces/%I
UMask=0077
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=yes
PrivateIPC=yes
PrivateNetwork=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
InaccessiblePaths=/var/lib/lga/audit /var/lib/lga/memory /var/lib/lga/uimp-spool
ReadOnlyPaths=/var/lib/lga/quarantine /etc/lga
ReadWritePaths=/var/lib/lga/jobs/%I /var/lib/lga/workspaces/%I /var/lib/lga/artifacts
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=~@mount @reboot @swap @obsolete @privileged
CPUAccounting=yes
CPUQuota=200%
IPAddressDeny=any
IOAccounting=yes
LimitFSIZE=16M
LimitNOFILE=1024
MemoryAccounting=yes
MemoryMax=4G
TasksAccounting=yes
TasksMax=512
TimeoutStartSec=30min
