# profile/airootfs/etc/systemd/system/getty@tty3.service.d/workspace-operator.conf
[Unit]
After=workspace-operator-setup.service
Requires=workspace-operator-setup.service

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin operator --noclear %I $TERM
